# submond.service
[Unit]
Description=submond service.

[Service]
Type=simple
ExecStart=/bin/bash /opt/submon/bin/submon.sh

[Install]
WantedBy=multi-user.target


